# kc-supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700                  ; socket file mode (default 0700)
chown=root:root             ; socket file uid:gid owner
[inet_http_server]
port=9000

[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug               ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


; ------------------------------------------------------------------

[program:tail]
command=/usr/bin/tail -f /dev/nul
priority=10
stdout_events_enabled=true
stderr_events_enabled=true